# deploy/supervisord.conf
# /etc/supervisor/conf.d/taskcaller.conf
# Runs the two app processes. Both self-load ~/live-caller/.env (python-dotenv /
# node dotenv), so no env wiring is needed here. HTTPS is handled by Caddy as its
# own systemd service (see deploy notes) — no tunnel process required.

[program:callserver]
command=/home/ubuntu/live-caller/.venv312/bin/python /home/ubuntu/live-caller/server.py
directory=/home/ubuntu/live-caller
environment=PYTHONUNBUFFERED=1
stdout_logfile=/var/log/callserver.log
stderr_logfile=/var/log/callserver.log
autostart=true
autorestart=true

[program:openwa]
command=node /home/ubuntu/live-caller/openwa/index.js
directory=/home/ubuntu/live-caller/openwa
stdout_logfile=/var/log/openwa.log
stderr_logfile=/var/log/openwa.log
autostart=true
autorestart=true
